# scripts/com.expense-tracker.streamlit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.expense-tracker.streamlit</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/expense-tracker/.venv/bin/streamlit</string>
        <string>run</string>
        <string>app/main.py</string>
        <string>--server.address</string>
        <string>192.168.3.10</string>
        <string>--server.port</string>
        <string>8501</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/expense-tracker</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DATA_DIR</key>
        <string>__HOME__/expense-data</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/expense-tracker.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/expense-tracker.log</string>
</dict>
</plist>
